# config/systemd/user/keyboard.service
[Unit]
Description=Update keyboard settings after startup

[Service]
Type=oneshot
ExecStart=/home/nicolai/dotfiles/scripts/keyboard.sh
