# config/systemd/user/keyboard.timer
[Unit]
Description=Update keyboard settings after startup

[Timer]
OnStartupSec=10seconds

[Install]
WantedBy=timers.target
